# boards/tinker/rootfs/etc/inittab
ttyS3::respawn:/sbin/getty -l /bin/bash -n 0 tty0
::sysinit:/etc/rc.init
::shutdown:/etc/rc.shutdown

